# instruments/nickel/vendor/lick_searchable_archive/deploy/roles/ingest_watchdog/templates/ingest_watchdog.service.j2
[Unit]
Description=Watchdog service to watch for new files and notify the archive of those files.
RequiresMountsFor={{ archive_data_mount }}

[Service]
Type=simple
ExecStart={{ venv_root }}/bin/python {{ venv_root }}/bin/ingest_watchdog.py --config {{ archive_config_dir }}/{{ watchdog_config }} -L {{ 'INFO' if archive_config == 'ops' else 'DEBUG' }} -l {{ archive_log_dir }}
User={{ archive_service_user }}
Group={{ archive_service_group }}
StandardOutput = null
StandardError= journal
# Allow group writeable files
UMask=0002
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
