# packages/lick_searchable_archive/deploy/roles/django_site/templates/gunicorn.service.j2
# Gunicorn systemd service config, adapted from https://docs.gunicorn.org/en/latest/deploy.html
[Unit]
Description=gunicorn daemon
After=network.target

[Service]
Type=notify
# the specific user that our service will run as
User={{ archive_service_user }}
Group={{ archive_service_group }}
# another option for an even more restricted service is
# DynamicUser=yes
# see http://0pointer.net/blog/dynamic-users-with-systemd.html
RuntimeDirectory=gunicorn
WorkingDirectory={{ venv_root }}/var/run
ExecStart={{ venv_root }}/bin/gunicorn -c {{ archive_config_dir }}/gunicorn.conf.py
ExecReload={{ venv_root }}/bin/kill -s HUP $MAINPID
# Allow group writeable files
UMask=0002
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
